# root/etc/systemd/system/checkupdates.timer
[Unit]
Description=Check available updates

[Timer]
OnCalendar=hourly
Unit=checkupdates.service

[Install]
WantedBy=timers.target
